# config/etc/cron.d/sks-dump-short-key-ids

0 0 * * * root timeout 12h /root/sks-dump-short-key-ids
