# config/etc/cron.d/sks-dump-short-key-ids
0 0 * * * root timeout 1h /root/sks-dump-short-key-ids && /root/sks-dump-short-key-ids-ping-healthcheck

