# config/etc/cron.d/sks-dump-short-key-ids

0 0 * * * root run-one /root/sks-dump-short-key-ids
